# surfingkeys/surfingkeys.service
[Unit]
Description=SurfingKeys Configuration Server
Documentation=https://github.com/brookhong/Surfingkeys
After=network.target
Wants=network.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/home/john/dotfiles-stow/surfingkeys
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=surfingkeys-server

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/john/dotfiles-stow/surfingkeys

# Environment
Environment=NODE_ENV=production
Environment=PORT=3030

[Install]
WantedBy=multi-user.target